# apps/render-bundle/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pass_environment=PORT,HOSTED_CHAIN_ID

[program:anvil]
command=/usr/local/bin/anvil --host 127.0.0.1 --port 8545 --chain-id 1337421 --block-time 3 --accounts 10 --balance 10000
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:gateway]
command=node /app/apps/render-bundle/scripts/gateway.mjs
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=ANVIL_RPC="http://127.0.0.1:8545"
